# rootfs/etc/services.d/chronyd/run
#!/bin/execlineb -P
/usr/sbin/chronyd -d -m -u chrony
